# scripts/nvchecker/nvchecker.service
[Unit]
Description=nvchecker update job

[Service]
User=pkger
Type=oneshot
ExecStart=/usr/bin/nvcheckerTimer

[Install]
WantedBy=multi-user.target
