# scripts/nvchecker/nvchecker.timer
[Unit]
Description=Run nvchecker every 5 minutes

[Timer]
OnUnitActiveSec=5m
Unit=nvchecker.service

[Install]
WantedBy=timers.target
